# scripts/systemd/backend.service
[Unit]
Description=Energy Platform Backend Server
After=network.target postgresql.service

[Service]
Type=simple
User=metal
WorkingDirectory=/home/metal/energy-platform/backend
Environment="PATH=/home/metal/energy-platform/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/metal/energy-platform/backend/venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


